# addon-cloudflare-tunnel/rootfs/etc/services.d/argo/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: addon-cloudflare-tunnel
# Starts Cloudflare Tunnel
# ==============================================================================
declare -a options
declare certificate 

certificate=$(bashio::config 'certificate')


bashio::log.info "Starting Cloudflare Tunnel"

# / Check URL
exec /opt/cloudflared --protocol quic --origincert=${certificate} tunnel --config=/config/cf-argo/config.yml run

Sleep 15
